# systemd/backup.service
[Unit]
Description=System and Home Backup Service
After=graphical-session.target
Wants=backup.timer

[Service]
Type=oneshot
# Replace USERNAME with your actual username
# Replace UID with your user ID (get with: id -u)
User=USERNAME
Group=USERNAME
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/UID
ExecStart=/usr/local/bin/backup-master.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
